# docker/supervisord-usps-processor.conf
[program:usps-processor]
command=/start-usps-processor.sh
numprocs=1
autostart=true
autorestart=true
