# docker/supervisord-usps-processor-importer.conf
[program:usps-processor-importer]
command=/start-usps-processor-importer.sh
numprocs=1
autostart=true
autorestart=true
