# docker/supervisord-usps-processor-api.conf
[program:usps-processor-api]
command=/start-usps-processor-api.sh
numprocs=1
autostart=true
autorestart=true
